# k0s.target
[Unit]
Description=k0s - Zero Friction Kubernetes
Documentation=https://docs.k0sproject.io

[Install]
WantedBy=multi-user.target
